# deploy/systemd/forward-dynatrace-handoff.service
[Unit]
Description=Forward for Dynatrace package handoff ingress
Documentation=https://github.com/forwardnetworks/forward-dynatrace
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=forward-dynatrace
Group=forward-dynatrace
WorkingDirectory=/opt/forward-dynatrace
EnvironmentFile=/etc/forward-dynatrace/forward-handoff.env
ExecStart=/usr/bin/node /opt/forward-dynatrace/scripts/forward-handoff-server.mjs
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/forward-dynatrace /var/log/forward-dynatrace
StateDirectory=forward-dynatrace
LogsDirectory=forward-dynatrace
UMask=0077

[Install]
WantedBy=multi-user.target
